# supervisor/supervisord.conf
[unix_http_server]
file=/root/data/projects/github.com/datacollector/supervisor/supervisor.sock

[supervisord]
logfile=/root/data/projects/github.com/datacollector/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/root/data/projects/github.com/datacollector/supervisor/supervisord.pid
nodaemon=false
silent=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///root/data/projects/github.com/datacollector/supervisor/supervisor.sock

[program:shrimp-daemon]
command=/root/data/projects/github.com/datacollector/venv/bin/python3 scheduler/daemon.py
directory=/root/data/projects/github.com/datacollector
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/root/data/projects/github.com/datacollector/logs/scheduler.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
