# supervisor/programs/scheduler.conf
[program:datacollector-scheduler]
command=/root/data/projects/github.com/datacollector/venv/bin/python3 scheduler/daemon.py
directory=/root/data/projects/github.com/datacollector
user=root
autostart=true
autorestart=true
startretries=5
startsecs=3
stderr_logfile=/root/data/projects/github.com/datacollector/logs/scheduler.log
stdout_logfile=/root/data/projects/github.com/datacollector/logs/scheduler.log
priority=25
